# systemd/spdyn.timer
[Unit]
Description=Submit IP to SPDyn

[Timer]
OnCalendar=0/3:00:00
Persistent=true

[Install]
WantedBy=timers.target
